# deploy/bridge.service
[Unit]
Description=TrainTracks ICS/SCADA Bridge Lab (static web server)
Documentation=file:///var/www/html/README.md
After=network.target

[Service]
Type=simple
# WEB_ROOT and PORT are substituted by deploy.sh; edit here if deploying by hand.
WorkingDirectory=/var/www/html
ExecStart=/usr/bin/python3 -m http.server --bind 0.0.0.0 8080 
Restart=on-failure
RestartSec=3
User=www-data
Group=www-data
# Light hardening — fine for an offline static file server.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
